# whatnow-backend.service
[Unit]
Description=WhatNow FastAPI backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/whatnow-backend
EnvironmentFile=/home/admin/whatnow-backend/.env
ExecStart=/home/admin/whatnow-backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8020
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
